# wttr.timer
[Unit]
Description=Run wttr script hourly

[Timer]
Unit=wttr.service
OnCalendar=minutely
Persistent=true

[Install]
WantedBy=timers.target
